# archiso/airootfs/root/config/mnt-mobihome.automount
[Unit]
Description=Automount unit for mobihome

[Automount]
Where=/mnt/mobihome

[Install]
WantedBy=remote-fs.target
