# archiso/airootfs/root/config/mnt-mobihome.mount
[Unit]
Description=Mount point for mobihome

[Mount]
What=mobihome
Where=/mnt/mobihome
Type=vboxsf
Options=uid=master,gid=master
